# metasploitable3/cookbooks/metasploitable/files/log4shell/log4shell.conf
description 'Run log4shell.conf App'
author 'jjsm'

start on runlevel [2345]

script
    sudo -u chewbacca /opt/log4shell/start.sh
end script
